# .config/systemd/user/vdirsyncer.timer
[Unit]
Description=vdirsyncer timer

[Timer]
OnCalendar=*:0/10
Unit=vdirsyncer.service

[Install]
WantedBy=timers.target
